# backend/cmake/cpack/postinst.in
#!/bin/sh
set -e
case "$1" in
configure|abort-upgrade|abort-deconfigure|abort-remove)
  ;;
*)
  exit 0
  ;;
esac
if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload >/dev/null 2>&1 || true
  systemctl enable ros_gui_backend.service >/dev/null 2>&1 || true
  systemctl start ros_gui_backend.service >/dev/null 2>&1 || true
fi
exit 0
